# public/vps-scraper-template/bca-scraper.service
[Unit]
Description=BCA Scraper v4.0.0 - Ultra-Robust Mode
After=network-online.target
Wants=network-online.target
Requires=network-online.target

# Restart rate limiting - max 5 restarts in 5 minutes
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/root/bca-scraper

# Pre-start checks
ExecStartPre=/bin/bash -c 'echo "=== BCA Scraper v4.0.0 Starting ===" >> /var/log/bca-scraper.log'
ExecStartPre=/bin/bash -c 'echo "Timestamp: $(date -Iseconds)" >> /var/log/bca-scraper.log'
ExecStartPre=/bin/bash -c 'test -f /root/bca-scraper/config.env || (echo "ERROR: config.env not found" >> /var/log/bca-scraper.log && exit 1)'
ExecStartPre=/bin/bash -c 'test -f /root/bca-scraper/bca-scraper.js || (echo "ERROR: bca-scraper.js not found" >> /var/log/bca-scraper.log && exit 1)'

# Kill any orphan chromium processes before starting
ExecStartPre=/bin/bash -c 'pkill -f "chromium.*puppeteer" 2>/dev/null || true'

# Wait for network
ExecStartPre=/bin/bash -c 'sleep 3'

# Main process - Persistent Browser Mode (single daemon)
ExecStart=/usr/bin/node /root/bca-scraper/bca-scraper.js

# Post-stop logging
ExecStopPost=/bin/bash -c 'echo "BCA Scraper stopped with exit code: $EXIT_STATUS at $(date -Iseconds)" >> /var/log/bca-scraper.log'

# Kill orphan processes on stop
ExecStopPost=/bin/bash -c 'pkill -f "chromium.*puppeteer" 2>/dev/null || true'

# Restart configuration
Restart=always
RestartSec=30

# Watchdog - systemd will restart if no response in 10 minutes
WatchdogSec=600

# No resource limits - high-spec VPS

# Logging
StandardOutput=append:/var/log/bca-scraper.log
StandardError=append:/var/log/bca-scraper-error.log
SyslogIdentifier=bca-scraper

# Environment
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=4096
Environment=TZ=Asia/Jakarta
EnvironmentFile=/root/bca-scraper/config.env

# Security - Chromium needs full system access
PrivateTmp=false
ProtectSystem=false

# Timeouts
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
