# public/vps-scraper-template/bca-scraper.service
[Unit]
Description=BCA Scraper v4.0.0 - Ultra-Robust Mode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/bca-scraper

# Pre-start: kill orphan chromium (- prefix ignores errors)
ExecStartPre=-/usr/bin/pkill -f "chromium.*puppeteer"
ExecStartPre=/bin/sleep 3

# Main process - Persistent Browser Mode
ExecStart=/usr/bin/node /root/bca-scraper/bca-scraper.js

# Post-stop cleanup (- prefix ignores errors)
ExecStopPost=-/usr/bin/pkill -f "chromium.*puppeteer"

# Auto-restart configuration
Restart=always
RestartSec=30

# Logging
StandardOutput=append:/var/log/bca-scraper.log
StandardError=append:/var/log/bca-scraper-error.log
SyslogIdentifier=bca-scraper

# Environment
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=4096
Environment=TZ=Asia/Jakarta
EnvironmentFile=-/root/bca-scraper/config.env

# Security - Chromium needs full system access
PrivateTmp=false
ProtectSystem=false

# Timeouts
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
