# public/vps-scraper-template/linuxsaldo/bca-saldo.service
# ============================================================
# BCA Balance Checker - Linux VPS v1.0.0
# Systemd Service File
# ============================================================
# Copy to: /etc/systemd/system/bca-saldo.service
# ============================================================

[Unit]
Description=BCA Balance Checker Linux v1.0.0
Documentation=https://github.com/your-repo
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bca-saldo
ExecStartPre=-/usr/bin/pkill -f "bca-balance-linux"
ExecStart=/usr/bin/node bca-balance-linux.js
ExecStopPost=-/usr/bin/pkill -9 -f chromium
Restart=always
RestartSec=10

# Environment
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bca-saldo

# Resource limits
MemoryMax=1G
CPUQuota=80%

# Security
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
